# github_upload/deploy/sput.service
[Unit]
Description=SPUT Flask app
After=network.target mysql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/sput
EnvironmentFile=/var/www/sput/.env
ExecStartPre=/var/www/sput/.venv/bin/python /var/www/sput/init_db.py
ExecStart=/var/www/sput/.venv/bin/gunicorn -c /var/www/sput/deploy/gunicorn.conf.py wsgi:application
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
